# supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /django/uwsgi.ini
autostart=true
autorestart=true

[program:app-nginx]
command = /usr/sbin/nginx
autostart=true
autorestart=true